# config/supervisor/conf.d/master-services.conf
[program:namenode]
command=/opt/hadoop/bin/hdfs --config /opt/hadoop/etc/hadoop namenode
user=root
autostart=true
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/opt/hadoop/logs/namenode.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hadoop/logs/namenode-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=10

[program:resourcemanager]
command=/opt/hadoop/bin/yarn --config /opt/hadoop/etc/hadoop resourcemanager
user=root
autostart=true
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/opt/hadoop/logs/resourcemanager.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hadoop/logs/resourcemanager-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=20

[program:hbase-master]
command=/opt/hbase/bin/hbase --config /opt/hbase/conf master start
user=root
autostart=true
autorestart=true
startsecs=30
startretries=5
stdout_logfile=/opt/hbase/logs/hbase-master.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hbase/logs/hbase-master-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=30

[program:hive-metastore]
command=/opt/hive/bin/hive --service metastore
user=root
autostart=true
autorestart=true
startsecs=30
startretries=10
stdout_logfile=/opt/hive/logs/hive-metastore.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hive/logs/hive-metastore-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=40

[program:hive-server2]
command=/opt/hive/bin/hive --service hiveserver2
user=root
autostart=true
autorestart=true
startsecs=30
startretries=10
stdout_logfile=/opt/hive/logs/hive-server2.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hive/logs/hive-server2-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=50

[program:spark-master]
command=/opt/spark/bin/spark-class org.apache.spark.deploy.master.Master
user=root
autostart=true
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/opt/spark/logs/spark-master.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/spark/logs/spark-master-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=60

[program:flink-jobmanager]
command=/opt/flink/bin/jobmanager.sh start-foreground
user=root
environment=FLINK_LOG_PREFIX="/opt/flink/log/flink--standalonesession"
autostart=true
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/opt/flink/log/flink-jobmanager.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/flink/log/flink-jobmanager-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=70

[program:health-monitor]
command=/opt/scripts/health-monitor.sh
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/health-monitor-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/health-monitor-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
priority=99
