# config/all-in-one/supervisor/infra-services.conf
[program:mysql]
command=/usr/bin/mysqld_safe
user=root
autostart=false
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/var/log/mysql/mysql.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/mysql/mysql-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=5
stopasgroup=true
killasgroup=true

[program:flume-agent]
command=/opt/flume/bin/flume-ng agent --conf /opt/flume/conf --conf-file /opt/flume/conf/flume-kafka.conf --name agent1 -Dflume.root.logger=INFO,console
user=root
autostart=false
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/opt/flume/logs/flume.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/flume/logs/flume-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=20

[program:health-monitor]
command=/opt/scripts/health-monitor.sh
user=root
autostart=false
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/health-monitor-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/health-monitor-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
priority=99
