# config/all-in-one/supervisor/worker-services.conf
[program:datanode]
command=/opt/hadoop/bin/hdfs --config /opt/hadoop/etc/hadoop datanode
user=root
autostart=false
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/opt/hadoop/logs/datanode.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hadoop/logs/datanode-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=10

[program:nodemanager]
command=/opt/hadoop/bin/yarn --config /opt/hadoop/etc/hadoop nodemanager
user=root
autostart=false
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/opt/hadoop/logs/nodemanager.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hadoop/logs/nodemanager-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=20

[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
user=root
autostart=false
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/opt/zookeeper/logs/zookeeper.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/zookeeper/logs/zookeeper-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=5

[program:hbase-regionserver]
command=/opt/scripts/hbase-foreground.sh regionserver
user=root
autostart=false
autorestart=true
startsecs=30
startretries=5
stdout_logfile=/opt/hbase/logs/hbase-regionserver.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/hbase/logs/hbase-regionserver-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=30
stopasgroup=true
killasgroup=true

[program:kafka-broker]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
user=root
autostart=false
autorestart=true
startsecs=20
startretries=5
stdout_logfile=/opt/kafka/logs/kafka.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/kafka/logs/kafka-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=25

[program:spark-worker]
command=/opt/spark/bin/spark-class org.apache.spark.deploy.worker.Worker spark://master:7077
user=root
autostart=false
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/opt/spark/logs/spark-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/spark/logs/spark-worker-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=40

[program:flink-taskmanager]
command=/opt/flink/bin/taskmanager.sh start-foreground
user=root
environment=HADOOP_CLASSPATH="/opt/hadoop/etc/hadoop:/opt/hadoop/share/hadoop/common/lib/*:/opt/hadoop/share/hadoop/common/*:/opt/hadoop/share/hadoop/hdfs/lib/*:/opt/hadoop/share/hadoop/hdfs/*:/opt/hadoop/share/hadoop/mapreduce/lib/*:/opt/hadoop/share/hadoop/mapreduce/*:/opt/hadoop/share/hadoop/yarn/lib/*:/opt/hadoop/share/hadoop/yarn/*"
autostart=false
autorestart=true
startsecs=15
startretries=5
stdout_logfile=/opt/flink/log/flink-taskmanager.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/opt/flink/log/flink-taskmanager-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=50

[program:health-monitor]
command=/opt/scripts/health-monitor.sh
user=root
autostart=false
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/health-monitor-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/health-monitor-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
priority=99
